# deploy/supervisor/modelizm-horizon.conf
; Supervisor: очереди Laravel Horizon (медиа-конверсии, письма, вебхуки оплат, ИИ-черновики).
; Путь: /etc/supervisor/conf.d/modelizm-horizon.conf
; Применить: sudo supervisorctl reread && sudo supervisorctl update

[program:modelizm-horizon]
process_name=%(program_name)s
command=php /var/www/modelizmclub/current/artisan horizon
directory=/var/www/modelizmclub/current
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/modelizm-horizon.log
stopwaitsecs=3600
; graceful shutdown — Horizon сам завершает воркеры по сигналу TERM
stopsignal=TERM
