# deploy/supervisor/modelizm-scheduler.conf
; Supervisor: планировщик Laravel (бэкапы БД, очистка кодов, истечение подписок).
; Альтернатива cron'у: держим sch:work постоянно запущенным.
; Путь: /etc/supervisor/conf.d/modelizm-scheduler.conf

[program:modelizm-scheduler]
process_name=%(program_name)s
command=php /var/www/modelizmclub/current/artisan schedule:work
directory=/var/www/modelizmclub/current
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/modelizm-scheduler.log
