# deploy/supervisor/modelizm-reverb.conf
; Supervisor: WebSocket-сервер Laravel Reverb (realtime чат и уведомления).
; Путь: /etc/supervisor/conf.d/modelizm-reverb.conf
; Применить: sudo supervisorctl reread && sudo supervisorctl update

[program:modelizm-reverb]
process_name=%(program_name)s
command=php /var/www/modelizmclub/current/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/modelizmclub/current
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/modelizm-reverb.log
stopwaitsecs=10
